# systemd/ultron.service.example
[Unit]
Description=Ultron (Discord ↔ Redmine ↔ LLM bot)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Repos/ultron-redmine
EnvironmentFile=/root/Repos/ultron-redmine/.env
Environment=CONFIG_PATH=/root/Repos/ultron-redmine/config.yaml
Environment=ULTRON_STATE_DIR=/root/Repos/ultron-redmine/data
# Include ~/.local/bin when cursor-agent is installed for the service user (Amvara /ca fallback).
Environment=PATH=/root/Repos/ultron-redmine/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root/.local/bin
ExecStart=/root/Repos/ultron-redmine/.venv/bin/python -m ultron
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillMode=mixed
StandardOutput=append:/root/Repos/ultron-redmine/ultron.log
StandardError=append:/root/Repos/ultron-redmine/ultron.log

[Install]
WantedBy=multi-user.target
